# systemd/etc_systemd_system_agent-platform-health.service
[Unit]
Description=AI Agent Platform Health Check

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/agent-platform-health
