# systemd/etc_systemd_system_backup.timer
[Unit]
Description=Daily backup timer (local + off-site)

[Timer]
# 03:30 local — runs after docker-prune.timer (Sun 03:00) so we capture
# the post-prune state and don't fight it for I/O.
OnCalendar=*-*-* 03:30:00
Persistent=true
RandomizedDelaySec=15min
Unit=backup.service

[Install]
WantedBy=timers.target
