# systemd/etc_systemd_system_agent-platform-health.timer
[Unit]
Description=Run platform health check 2 minutes after boot and then hourly

[Timer]
OnBootSec=2min
OnUnitActiveSec=1h
Unit=agent-platform-health.service

[Install]
WantedBy=timers.target
